# toolboxes/cli/files/etc/init.d/nix-daemon
#!/sbin/openrc-run
description="Nix multi-user support daemon"

command="/nix/var/nix/profiles/default/bin/nix-daemon"
command_background="yes"
command_args_background="--daemon"
pidfile="/run/$RC_SVCNAME.pid"

start() {
        ebegin "Starting ${SVCNAME}"
        start-stop-daemon --start --quiet --background --exec $command
        eend ${?}
}

stop() {
        ebegin "Stopping ${SVCNAME}"
        start-stop-daemon --stop --quiet --exec $command
        eend ${?}
}
